# services/kinto/kinto.container
[Container]
Image=docker.io/kinto/kinto-server:latest
